# conf/etc/service/redis/run
#!/bin/bash

[[ -n ${REDIS_TAGDB} ]] || exit 1
exec /usr/bin/redis-server /etc/redis/redis.conf