# conf/etc/service/carbon-aggregator/run
#!/bin/sh

[[ -f "./down" ]] && exit 1

BIN=python3
[[ -f "/opt/graphite/bin/pypy3" ]] && BIN=pypy
exec ${BIN} /opt/graphite/bin/carbon-aggregator.py start --debug 2>&1
