# conf/etc/service/nginx/run
#!/bin/bash

mkdir -p /var/log/nginx
/usr/sbin/nginx -c /etc/nginx/nginx.conf
