# conf/etc/service/graphite/run
#!/bin/bash

exec /usr/bin/python /usr/local/src/graphite-web/webapp/manage.py runfcgi daemonize=false host=127.0.0.1 port=8080

