# conf/etc/service/cron/run
#!/bin/sh

[[ -f "./down" ]] && exit 1

exec /usr/sbin/crond -f
