# conf/etc/service/grafana/run
#!/bin/bash

exec /usr/sbin/grafana-server -config /etc/grafana/grafana.ini -homepath /usr/share/grafana
