# conf/etc/service/collectd/run
#!/bin/sh

[ -n "${COLLECTD}" ] || exit 1
exec /usr/sbin/collectd -f -C /etc/collectd/collectd.conf
